# bot.conf
[program:lawyerbot]
command=/home/bot/lawyerbot/run_bot.sh
autostart=true
autorestart=true
stderr_logfile=/var/log/lawyerbot.err.log
stdout_logfile=/var/log/lawyerbot.out.log
